# files/etc/periodic/daily/backup.sh
#!/bin/bash

ruby /backup_script.rb run
